# specs/default/chef/site-cookbooks/gridengine/templates/default/sgeexecd.service.erb
[Unit]
Description=Univa Grid Engine Execd
After=network.target remote-fs.target

[Service]
Type=forking
KillMode=process
ExecStart=<%= @gridengineroot %>/<%= @gridenginecell %>/common/sgeexecd start
ExecStop=<%= @gridengineroot %>/<%= @gridenginecell %>/common/sgeexecd stop

[Install]
WantedBy=multi-user.target